# crontab/asf-sync
3 5,17 * * * * ubuntu bash /home/ubuntu/asf-mirror/sync_mirror.sh