# crontab/asf-sync
5 0,12 * * * ubuntu bash /home/ubuntu/asf-mirror/sync_mirror.sh